# zfix-docker@.service.in
[Unit]
Description=zfix.org dockerized service '%I'
After=docker.service
Requires=docker.service

[Service]
WorkingDirectory={{REPO_PATH}}/applications/%i/
# ExecStartPre=-/usr/local/bin/docker-compose pull --ignore-pull-failures
# ExecStartPre=-/usr/local/bin/docker-compose build
ExecStart=/usr/local/bin/docker-compose up

[Install]
WantedBy=multi-user.target
