# zfix-docker@.service.in
[Unit]
Description=zfix.org dockerized service '%I'
After=docker.service
Requires=docker.service

[Service]
# ExecStartPre=-/usr/local/bin/docker-compose -f {{REPO_PATH}}/applications/%i/docker-compose.yml pull --ignore-pull-failures
# ExecStartPre=-/usr/local/bin/docker-compose -f {{REPO_PATH}}/applications/%i/docker-compose.yml build
ExecStart=/usr/local/bin/docker-compose -f {{REPO_PATH}}/applications/%i/docker-compose.yml up

[Install]
WantedBy=multi-user.target