# zfix-docker@.service.in
[Unit]
Description=zfix.org docker application service
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=-/usr/local/bin/docker-compose -f {{REPO_PATH}}/applications/%i/docker-compose.yml build --pull
ExecStart=/usr/local/bin/docker-compose -f {{REPO_PATH}}/applications/%i/docker-compose.yml up

[Install]
WantedBy=multi-user.target